# JPG/splash_anim
#!/bin/sh /etc/rc.common
START=10
STOP=99

# --- 动态探测函数 ---
find_fb_device() {
    for fb in /sys/class/graphics/fb*; do
        [ -e "$fb/name" ] || continue
        name=$(cat "$fb/name")
        if echo "$name" | grep -qiE "mipi|spi|panel-mipi-dbi"; then
            echo "/dev/$(basename $fb)"
            return 0
        fi
    done
    for fb in /sys/class/graphics/fb*; do
        [ -e "$fb/dev" ] && echo "/dev/$(basename $fb)" && return 0
    done
    return 1
}

find_tty_number() {
    for i in 1 2 3 4 5 6; do
        if [ -c "/dev/tty$i" ]; then
            echo "$i"
            return 0
        fi
    done
    [ -c "/dev/tty0" ] && echo "0" && return 0
    return 1
}

start() {
    # 等待fb就绪
    local i=0
    FB_DEV=""
    while [ $i -lt 20 ]; do
        FB_DEV=$(find_fb_device)
        [ -n "$FB_DEV" ] && break
        usleep 500000
        i=$((i + 1))
    done

    if [ -z "$FB_DEV" ]; then
        logger -t splash "ERROR: No framebuffer found, exit splash"
        return 1
    fi

    # 点亮背光（仅支持0/1）
    BL_DIR="/sys/class/backlight/backlight"
    if [ -d "$BL_DIR" ]; then
        echo 1 > ${BL_DIR}/brightness 2>/dev/null
    fi

    logger -t splash "Start splash anim, fb=$FB_DEV"

    # 后台运行Python循环播放LOGO，替代fbi
    python3 /usr/bin/splash_loop.py &
    echo $! > /var/run/splash.pid
}

stop() {
    # 杀掉动画Python进程
    if [ -f /var/run/splash.pid ]; then
        kill $(cat /var/run/splash.pid) 2>/dev/null
        rm -f /var/run/splash.pid
    fi
    # 无fbi，不需要killall fbi

    # 重新探测fb，渲染网络短句画面
    local FB_DEV=$(find_fb_device)
    if [ -z "$FB_DEV" ]; then
        logger -t splash "stop stage: can not find fb device, skip text render"
        return 1
    fi

    sleep 2
    python3 /usr/bin/show_sentence.py
    # 直接用PIL写入fb，不再调用fbi
    python3 -c "from PIL import Image; img=Image.open('/tmp/screen_out.jpg').convert('RGB'); img.save('/dev/fb0','raw')" &
}
